# waittime.cron
*/10 * * * * python /Users/brlamore/src/waittime/waittimes.py >> /Users/brlamore/src/waittime/waittimes.log
